# .config/systemd/system/dropbox.service
[Unit]
Description=Dropbox Daemon
After=network.target

[Service]
ExecStart=/home/gph/.dropbox-dist/dropboxd

[Install]
WantedBy=multi-user.target
